# skills/devops/cursor-subapi-compat/scripts/cursor-cpa-compat.service
[Unit]
Description=Cursor compatibility bridge for direct CPA
After=network.target cliproxyapi.service

[Service]
Type=simple
WorkingDirectory=/root/cursor-cpa-compat
ExecStart=/usr/bin/python3 /root/cursor-cpa-compat/server.py
Restart=always
RestartSec=2
User=root
Environment=CURSOR_CPA_DEBUG_SSE=0

[Install]
WantedBy=multi-user.target
